# systemd/glycan-analysis-pipeline.service
[Unit]
Description=Run GlycanAnalysisPipeline at boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/scratch/os1e24/GlycanAnalysisPipeline
Environment=HOME=/scratch/os1e24
Environment=PATH=/scratch/os1e24/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/scratch/os1e24/GlycanAnalysisPipeline/boot-run.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
